# systemd/ssh-noti-summary.timer
[Unit]
Description=SSH Notifier Summary (hourly)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
